# calibre-web-automated/root/etc/s6-overlay/s6-rc.d/svc-calibre-web-automated/run
#!/usr/bin/with-contenv bash

set -euo pipefail

python3 - <<'PY' >/tmp/cwa-addon-env
import json
from pathlib import Path

options = {}
options_path = Path('/data/options.json')
if options_path.exists():
    options = json.loads(options_path.read_text())

for key, default in {
    'TZ': 'UTC',
    'PUID': '0',
    'PGID': '0',
    'CONFIG_PATH': '/addon-config',
    'INGEST_PATH': '/share/calibre-web-automated/ingest',
    'LIBRARY_PATH': '/media/calibre-web-automated/library',
    'NETWORK_SHARE_MODE': 'false',
    'HARDCOVER_TOKEN': '',
    'TRUSTED_PROXY_COUNT': '1',
    'DISABLE_LIBRARY_AUTOMOUNT': 'false',
    'CWA_PORT_OVERRIDE': '8083',
}.items():
    value = options.get(key, default)
    if isinstance(value, bool):
        value = 'true' if value else 'false'
    print(f'{key}={value}')
PY

while IFS='=' read -r key value; do
    export "$key=$value"
done < /tmp/cwa-addon-env

rm -f /tmp/cwa-addon-env

CONFIG_PATH="${CONFIG_PATH:-/addon-config}"
INGEST_PATH="${INGEST_PATH:-/share/calibre-web-automated/ingest}"
LIBRARY_PATH="${LIBRARY_PATH:-/media/calibre-web-automated/library}"

exec s6-setuidgid abc python3 /app/calibre-web-automated/cps.py